# hyprmode-daemon.service
[Unit]
Description=HyprMode Emergency Recovery Daemon
After=graphical-session.target

[Service]
Type=simple
ExecStartPre=/bin/bash -c 'find /usr/local/bin -name "*.pyc" -delete 2>/dev/null; find ${HOME}/.cache -name "*.pyc" -delete 2>/dev/null || true'
ExecStart=/usr/local/bin/hyprmode-daemon-wrapper
Restart=on-failure
RestartSec=5

[Install]
WantedBy=default.target
