# hyprmode-daemon.service
[Unit]
Description=HyprMode Emergency Recovery Daemon
After=graphical-session.target

[Service]
Type=simple
ExecStart=/usr/local/bin/hyprmode-daemon
Restart=on-failure
RestartSec=5

[Install]
WantedBy=default.target

